# deploy/systemd/sales-marketplaces-backend.service
[Unit]
Description=Sales Marketplaces Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sales-marketplaces
EnvironmentFile=/opt/sales-marketplaces/.env
ExecStart=/opt/sales-marketplaces/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 18000
Restart=always
RestartSec=5
TimeoutStopSec=15
KillMode=mixed
User=root

[Install]
WantedBy=multi-user.target
